# deploy/binance-testnet-report.service
[Unit]
Description=Send Binance Spot Live Agent daily email report

[Service]
Type=oneshot
WorkingDirectory=/opt/binance-testnet-agent
EnvironmentFile=/opt/binance-testnet-agent/.env
ExecStart=/usr/bin/python3 -m binance_testnet_agent.email_report --period day
User=ubuntu
Group=ubuntu
